# contrib/hoard@.service
# Hoard systemd template unit
# Usage: systemctl enable --now hoard@myapp
#
# Creates a dedicated Hoard instance bound to %i.service.
# The control socket will be at /run/hoard/%i.sock
#
# Requires: hoard binary at /usr/local/bin/hoard
#           BPF object at /usr/lib/hoard/hoard.bpf.o

[Unit]
Description=Hoard backup daemon for %i
Documentation=https://github.com/hoard-project/hoard
After=network-online.target
BindsTo=%i.service
Wants=network-online.target

[Service]
Type=notify
ExecStart=/usr/local/bin/hoard \
  --mode standalone \
  --service %i \
  --watch-root /var/lib/%i/data \
  --tls-mode ktls \
  --s3-endpoint ${HOARD_S3_ENDPOINT} \
  --s3-region ${HOARD_S3_REGION} \
  --s3-bucket ${HOARD_S3_BUCKET} \
  --s3-access-key ${HOARD_S3_ACCESS_KEY} \
  --s3-secret-key ${HOARD_S3_SECRET_KEY} \
  --s3-prefix hoard/%i \
  --gc-interval 21600 \
  --gc-ttl-days 30
ExecReload=/bin/kill -HUP $MAINPID
KillSignal=SIGTERM
TimeoutStopSec=30

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/%i/data /run/hoard
PrivateTmp=yes
ProtectKernelTunables=no  # BPF loading needs kernel access
ProtectKernelModules=yes
ProtectControlGroups=yes
MemoryDenyWriteExecute=yes
RestrictRealtime=yes
RestrictNamespaces=yes
LockPersonality=yes

# BPF capabilities
CapabilityBoundingSet=CAP_BPF CAP_SYS_ADMIN CAP_NET_RAW
AmbientCapabilities=CAP_BPF CAP_SYS_ADMIN CAP_NET_RAW

# Resource limits
LimitNOFILE=65536
LimitMEMLOCK=infinity
CPUQuota=5%
MemoryMax=128M

[Install]
WantedBy=multi-user.target